# files/lxc-dnsmasq.conf
description "lxc dnsmasq service"
author "Adam Ryczkowski, ispired by Serge Hallyn <serge.hallyn@canonical.com>"

expect fork

start on started lxc-net
stop on stopped lxc-net

env USE_LXC_BRIDGE="true"
env LXC_NETWORK_CONFIGURATION_FILE="/etc/lxc/bridge-networks"
env LXC_DHCP_CONFILE="/etc/lxc/dnsmasq.conf"
env varrun="/run/lxc-dnsmasq"

pre-start script
	[ -f /etc/default/lxc ] && . /etc/default/lxc

	[ -f /etc/default/lxc-net ] && . /etc/default/lxc-net

	[ "x$USE_LXC_BRIDGE" = "xtrue" ] || { stop; exit 0; }

	if [ ! -d ${varrun} ]; then
		mkdir -p ${varrun}
	fi
	bash /usr/local/lib/lxc-scripts/lxc-bridge-parser.sh --conf-file $LXC_NETWORK_CONFIGURATION_FILE --remove-dhcp-leases

	opts="$LXC_DOMAIN_ARG --pid-file=${varrun}/dnsmasq.pid --conf-file=${LXC_DHCP_CONFILE} --dhcp-lease-max=${LXC_DHCP_MAX} --keep-in-foreground"

	echo "Executing:\ndnsmasq $opts"

	/usr/sbin/dnsmasq $opts &
end script

post-stop script
	if [ -f ${varrun}/dnsmasq.pid ]; then
		PID=`cat ${varrun}/dnsmasq.pid`
		kill $PID
	fi
end script